# com.carlos.clap-spotify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.carlos.clap-spotify</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/carlosoliva/clap-spotify/.venv/bin/python</string>
        <string>/Users/carlosoliva/clap-spotify/clap_play.py</string>
    </array>

    <!-- Logs sin buffer para poder depurar -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Arrancar al iniciar sesión -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Reiniciar si se cae -->
    <key>KeepAlive</key>
    <true/>

    <!-- No martillar el reinicio si falla seguido -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>WorkingDirectory</key>
    <string>/Users/carlosoliva/clap-spotify</string>

    <key>StandardOutPath</key>
    <string>/Users/carlosoliva/clap-spotify/clap.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/carlosoliva/clap-spotify/clap.err</string>
</dict>
</plist>
